# infra/ansible/roles/promtail/files/promtail.service
[Unit]
Description=Promtail — Loki log shipper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=promtail
Group=promtail
ExecStart=/usr/local/bin/promtail -config.file=/etc/promtail/config.yml
Restart=on-failure
RestartSec=5
StateDirectory=promtail
LogsDirectory=promtail

# Sandboxing
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/promtail
SupplementaryGroups=systemd-journal

[Install]
WantedBy=multi-user.target
